# backend/deploy/vayupos-backend.service
[Unit]
Description=VayuPos Backend API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/vayupos/backend
EnvironmentFile=/home/ubuntu/vayupos/backend/.env
ExecStart=/home/ubuntu/vayupos/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vayupos-backend

[Install]
WantedBy=multi-user.target
